# debian/homcc-server.service
[Unit]
Description=The homcc-server service.

[Service]
ExecStart=/usr/bin/homcc-server

[Install]
WantedBy=multi-user.target
