# debian/homccd.service
[Unit]
Description=The homccd service.

[Service]
ExecStart=/usr/bin/homccd

[Install]
WantedBy=multi-user.target
